# layout/DEBIAN/postinst
#!/usr/bin/env bash

exec 3>&1 &>/dev/null

if [ "$DPKG_MAINTSCRIPT_ARCH" == "iphoneos-arm64" ]; then
    DPKG_ROOT="/var/jb"
fi

SERV_AGENT="$DPKG_ROOT/Library/LaunchDaemons/com.zerzvx.waifuvnc.plist"

if [ -f "$SERV_AGENT" ]; then
    launchctl unload "$SERV_AGENT" || true
else
    launchctl remove system/com.zerzvx.waifuvnc || true
fi

sleep 2

if [ -f "$SERV_AGENT" ]; then
    launchctl load -w -F "$SERV_AGENT" || true
fi

# Registrar el icono de la app Waifu VNC en el home screen
if [ -d "$DPKG_ROOT/Applications/WaifuVNC.app" ]; then
    uicache -p "$DPKG_ROOT/Applications/WaifuVNC.app" || uicache -a || true
fi

exit 0